# recipes-tools/service-tools/files/watchdog.service
[Unit]
Description=WatchDog supervise

[Service]
Type=simple
ExecStart=/usr/sbin/watchdog.sh
Restart=always

[Install]
WantedBy=multi-user.target
